# deploy/systemd/pdf_parser.service
[Unit]
Description=Dockling API
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/html/pdf_parser/current
EnvironmentFile=-/var/www/html/pdf_parser/shared/.env
Environment=PDF_PARSER_TEMP_DIR=/tmp/pdf_parser
ExecStart=/var/www/html/pdf_parser/current/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 29999
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
